# operation-echo/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:ollama]
command=ollama serve
autostart=true
autorestart=true
priority=10
environment=OLLAMA_HOST="0.0.0.0:11434"
stdout_logfile=/var/log/ollama.log
stderr_logfile=/var/log/ollama.err

[program:app]
command=uvicorn main:app --host 127.0.0.1 --port 8000
directory=/opt/app
autostart=true
autorestart=true
priority=20
environment=OLLAMA_HOST="http://127.0.0.1:11434",MODEL_NAME="llama3.2:3b"
stdout_logfile=/var/log/app.log
stderr_logfile=/var/log/app.err

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.err
