# operation-nova/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:anvil]
command=/root/.foundry/bin/anvil --host 0.0.0.0 --port 8545 --balance 2000000
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/anvil.log
stderr_logfile=/var/log/anvil.err

# One-shot deploy job. Waits for anvil, deploys and seeds, then exits 0.
[program:deploy]
command=/opt/deploy.sh
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=20
environment=RPC_URL="http://127.0.0.1:8545",STATE_PATH="/opt/flagserver/state.json"
stdout_logfile=/var/log/deploy.log
stderr_logfile=/var/log/deploy.err

[program:flagserver]
command=python3 /opt/flagserver/app.py
directory=/opt/flagserver
autostart=true
autorestart=true
priority=30
environment=RPC_URL="http://127.0.0.1:8545",PORT="80",STATE_PATH="/opt/flagserver/state.json",FLAG_PATH="/opt/flagserver/flag.txt",CONTRACT_PATH="/opt/flagserver/NovaVault.sol"
stdout_logfile=/var/log/flagserver.log
stderr_logfile=/var/log/flagserver.err
